# Project2/part1/cafe/etc/rc.local
#!/bin/sh -e
#
# rc.local

iptables -F
iptables -X
iptables -P INPUT ACCEPT
iptables -P FORWARD ACCEPT
iptables -P OUTPUT ACCEPT
/etc/init.d/quagga start
exit 0
